# roles/https_haproxy/templates/ocsp.timer.j2
[Unit]
Description=Create OSCP timer

[Timer]
OnCalendar=01:00

[Install]
WantedBy=multi-user.target

