# roles/sram_monitor/templates/sram-monitor.timer.j2
[Unit]
Description=SRAM Monitor

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=multi-user.target
