# roles/sram_monitor/templates/sram-monitor.timer.j2
[Unit]
Description=SRAM Monitor

[Timer]
OnCalendar=*-*-* *:*/15:00

[Install]
WantedBy=multi-user.target
