# roles/sram_monitor/templates/sram-monitor@.timer.j2
[Unit]
Description=SRAM Monitor

[Timer]
OnCalendar={{item.timer}}

[Install]
WantedBy=multi-user.target
