# roles/slp/templates/scz-slp.timer.j2
[Unit]
Description=SCZ SLP Ldap Synchronization Timer

[Timer]
OnCalendar=*:00/5

[Install]
WantedBy=multi-user.target

