# roles/ldap_clients/templates/lsc.timer.j2
 [Unit]
Description=LSC timer

[Timer]
OnUnitActiveSec=300

[Install]
WantedBy=multi-user.target
