# roles/ldap_clients/templates/lsc.service.j2
[Unit]
Description=LSC Runner

[Service]
User=root
Group=root
Type=oneshot
WorkingDirectory={{ lsc_install_dir }}
ExecStart={{ lsc_install_dir }}/run_lsc_tasks

