# roles/backup-offsite/templates/sram-backup-sync.timer.j2
[Unit]
Description=SRAM Backup syncer

[Timer]
OnCalendar=*-*-* *:00/15:00

[Install]
WantedBy=multi-user.target
