# roles/ldap_monitor/templates/ldap_monitor.service.j2
[Unit]
Description=LDAP Checksum Runner

[Service]
User=root
Group=root
Type=oneshot
WorkingDirectory={{ ldap_checksum_dir }}
RuntimeMaxSec=120
ExecStart={{ ldap_monitor_dir }}/ldap_monitor.sh
