# roles/logging_collector/templates/sram-remotejournal.socket.j2
[Unit]
Description=Journal Remote Sink Socket

[Socket]
ListenStream={{item.value}}

[Install]
WantedBy=sockets.target
