# roles/aclsync/templates/sram_aclsync.timer.j2
[Unit]
Description=PLSC timer

[Timer]
OnBootSec=30
OnUnitActiveSec={{aclsync_interval_seconds}}

[Install]
WantedBy=multi-user.target
