# roles/scim_monitor/templates/sram-scimmonitor.timer.j2
[Unit]
Description=SCIM Monitor timer

[Timer]
OnCalendar=*:{{2+extra_min}}/15

[Install]
WantedBy=multi-user.target
