# roles/sbs/templates/sbs.service.j2
[Unit]
Description=SBS
After=syslog.target network.target

[Service]
Type=simple
User=sbs
WorkingDirectory={{ sbs_workdir }}
ExecStart=@/usr/local/bin/docker-compose docker-sbs up
Restart=never
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sbs

[Install]
WantedBy=multi-user.target
