# roles/backup_collector/templates/sram-processbackups.timer.j2
[Unit]
Description=SRAM Backup processor

[Timer]
OnCalendar=*-*-* *:00/10:00

[Install]
WantedBy=multi-user.target

