# roles/lb_haproxy/files/fence.service
[Unit]
Description=LDAP Fencing mechanism
After=haproxy.service

[Service]
ExecStart=/opt/fence.sh

[Install]
WantedBy=multi-user.target