# roles/tools_contacts/templates/sram_contacts.timer.j2
[Unit]
Description=SRAM contacts update

[Timer]
OnBootSec=30
OnUnitActiveSec={{contacts_interval_seconds}}

[Install]
WantedBy=multi-user.target
