# roles/ldap_checksum/templates/ldap_checksum.service.j2
[Unit]
Description=LDAP Checksum Runner

[Service]
User=root
Group=root
Type=oneshot
WorkingDirectory={{ ldap_checksum_dir }}
ExecStart={{ ldap_checksum_dir }}/checksum.sh
RuntimeMaxSec=120
